# deploy/systemd/polymarket-backup.timer
[Unit]
Description=Nightly Polymarket Data Backup

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
Unit=polymarket-backup.service

[Install]
WantedBy=timers.target
