# deploy/systemd/polymarket-disk-maintenance-staging.timer
[Unit]
Description=Hourly Ditto disk maintenance (staging)

[Timer]
OnCalendar=hourly
Persistent=true
Unit=polymarket-disk-maintenance-staging.service

[Install]
WantedBy=timers.target
