# deploy/systemd/polymarket-disk-maintenance.timer
[Unit]
Description=Hourly Ditto disk maintenance

[Timer]
OnCalendar=hourly
Persistent=true
Unit=polymarket-disk-maintenance.service

[Install]
WantedBy=timers.target
